# install/etc/services.d/dovecot/run
#!/usr/bin/with-contenv bash

while [ ! -f /tmp/state/dovecot-init ]; do
  sleep 1
done

/usr/sbin/dovecot -F
