# install/etc/services.d/dovecot/run
#!/usr/bin/with-contenv bash

echo ''
echo '** Starting Dovecot'
exec dovecot -F
